# templates/etc/systemd/system/paperless-web.service
[Unit]
Description=Paperless Web Daemon

[Service]
User=paperless
Group=paperless
WorkingDirectory=/opt/paperless/source/src
ExecStart={{ paperless_virtualenv }}/bin/gunicorn paperless.wsgi -w 2

[Install]
WantedBy=multi-user.target
